# deploy/call-hermes-ops.timer
[Unit]
Description=Run Call Hermes operational checks

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
