# deploy/call-hermes-backup.timer
[Unit]
Description=Daily Call Hermes database backup

[Timer]
OnCalendar=*-*-* 03:15:00
RandomizedDelaySec=15min
Persistent=true

[Install]
WantedBy=timers.target
