# deploy/call-hermes-backup.service.in
[Unit]
Description=Back up Call Hermes conversations

[Service]
Type=oneshot
User=@USER@
WorkingDirectory=@ROOT@
EnvironmentFile=@ROOT@/server/.env
ExecStart=@ROOT@/scripts/backup_conversations.sh
